# content-pipeline/launchd/com.ai5phut.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai5phut.bot</string>

    <!-- Use the wrapper script — it activates the venv so all deps (googleapiclient, etc.) are available.
         Replace /Users/YOU with your actual home directory path. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/ContentCreator/content-pipeline/run_pipeline.sh</string>
        <string>--bot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOU/ContentCreator/content-pipeline</string>

    <!-- Daemon: chạy liên tục, tự restart nếu crash -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOU/ContentCreator/content-pipeline/logs/bot_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/ContentCreator/content-pipeline/logs/bot_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Tự restart nếu crash, chờ 10s -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
