# content-pipeline/launchd/com.ai5phut.drama-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai5phut.drama-health</string>

    <!-- Chạy qua wrapper tĩnh run_module.sh (resolve venv ở runtime) thay vì
         trỏ thẳng venv/bin/python3 — tránh EX_CONFIG (78) khi venv bị rebuild
         (issue #74/#75; xem run_module.sh).
         Replace /Users/YOU with your actual home directory path. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/ContentCreator/content-pipeline/run_module.sh</string>
        <string>-m</string>
        <string>storage.collector_health</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOU/ContentCreator/content-pipeline</string>

    <!-- Chạy mỗi 12h (06:30 và 18:30) — alert Telegram nếu reddit_drama
         collector chưa chạy thành công quá 2 ngày (EPIC #2.4). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/YOU/ContentCreator/content-pipeline/logs/drama_health_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/ContentCreator/content-pipeline/logs/drama_health_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
