# content-pipeline/launchd/com.ai5phut.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai5phut.pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOU/ContentCreator/content-pipeline</string>

    <!-- Chạy lúc 7:00 sáng mỗi ngày -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOU/ContentCreator/content-pipeline/logs/pipeline_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/ContentCreator/content-pipeline/logs/pipeline_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
