# content-pipeline/launchd/com.ai5phut.weekly-retro.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai5phut.weekly-retro</string>

    <!-- Root cause #74/#75: plist KHÔNG khai báo WorkingDirectory/StandardOutPath/
         StandardErrorPath. launchd/xpcproxy dựng chúng TRƯỚC khi exec và giữ handle
         inode; khi thư mục bị xoá/tạo lại (rebuild venv, re-clone, reconfig) handle
         stale → EX_CONFIG (78) *trước khi* binary chạy, khoá job tới khi reload.
         run_module.sh thiết lập cwd + log ở runtime với inode tươi (LOG_BASENAME).
         Replace /Users/YOU with your actual home directory path. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/ContentCreator/content-pipeline/run_module.sh</string>
        <string>-m</string>
        <string>analytics.weekly_retro</string>
    </array>

    <!-- Chủ nhật 19:00 — đẩy báo cáo tuần qua Telegram để tối CN đọc, sáng T2
         đã có quyết định (Phase 6 — Weekly retro). Weekday: 0=CN trong launchd. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>LOG_BASENAME</key>
        <string>weekly_retro</string>
    </dict>
</dict>
</plist>
