# content-pipeline/systemd/ai5phut-pipeline.service
[Unit]
Description=AI 5 Phut Moi Ngay — Daily Content Pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/user/ContentCreator/content-pipeline
ExecStart=/usr/local/bin/python3 main.py
EnvironmentFile=/home/user/ContentCreator/content-pipeline/.env

# Logging — journalctl -u ai5phut-pipeline
StandardOutput=journal
StandardError=journal

# Timeout: pipeline can take ~10 min (AI calls + video gen)
TimeoutStartSec=900
