# content-pipeline/systemd/ai5phut-bot.service
[Unit]
Description=AI 5 Phut Moi Ngay — Telegram Bot (approve/publish)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/user/ContentCreator/content-pipeline
ExecStart=/usr/local/bin/python3 main.py --bot
EnvironmentFile=/home/user/ContentCreator/content-pipeline/.env
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
